# scripts/systemd/council-dawn-mist.service
[Unit]
Description=Council Dawn Mist Standup — Daily Council Review
After=network-online.target arxiv-crawler.service
Wants=network-online.target

[Service]
Type=oneshot
User=dereadi
WorkingDirectory=/ganuda
ExecStart=/home/dereadi/cherokee_venv/bin/python3 /ganuda/scripts/council_dawn_mist.py
Environment=HOME=/home/dereadi
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
