# scripts/systemd/credential-scanner.timer
[Unit]
Description=Credential Scanner — Weekly Saturday 2 AM

[Timer]
OnCalendar=Sat *-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
