# scripts/systemd/fire-guard.timer
[Unit]
Description=Fire Guard — check all federation services every 2 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=120
AccuracySec=15

[Install]
WantedBy=timers.target
