# scripts/systemd/moltbook-proxy.service
[Unit]
Description=Cherokee AI Federation - Moltbook Proxy (Crawdad)
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=dereadi
Group=dereadi
WorkingDirectory=/ganuda/services/moltbook_proxy
ExecStart=/home/dereadi/cherokee_venv/bin/python -u proxy_daemon.py
Restart=on-failure
RestartSec=60
StandardOutput=append:/ganuda/logs/moltbook_proxy.log
StandardError=append:/ganuda/logs/moltbook_proxy.log
Environment=PYTHONPATH=/ganuda/lib:/ganuda/services/moltbook_proxy
Environment=MOLTBOOK_API_URL=https://www.moltbook.com
Environment=CHEROKEE_DB_HOST=192.168.132.222
Environment=CHEROKEE_DB_NAME=zammad_production
Environment=CHEROKEE_DB_USER=claude
EnvironmentFile=/ganuda/config/secrets.env

[Install]
WantedBy=multi-user.target
