# scripts/systemd/safety-canary.service
[Unit]
Description=Safety Canary — Red team probe of LLM alignment
After=network-online.target llm-gateway.service

[Service]
Type=oneshot
User=dereadi
WorkingDirectory=/ganuda
ExecStart=/home/dereadi/cherokee_venv/bin/python3 /ganuda/scripts/safety_canary.py
Environment=HOME=/home/dereadi
Environment=PYTHONPATH=/ganuda:/ganuda/lib
EnvironmentFile=/ganuda/config/secrets.env
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
