# scripts/systemd/sag-unified.service
[Unit]
Description=SAG Unified Interface - Cherokee AI Control Room
After=network.target postgresql.service
Wants=visual-kanban.service

[Service]
Type=simple
User=dereadi
Group=dereadi
WorkingDirectory=/ganuda/home/dereadi/sag_unified_interface
Environment="PATH=/ganuda/home/dereadi/sag_unified_interface/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/ganuda/lib:/ganuda/home/dereadi/sag_unified_interface"
ExecStart=/ganuda/home/dereadi/sag_unified_interface/venv/bin/python3 app.py
Restart=always
RestartSec=10
StandardOutput=append:/ganuda/logs/sag-unified.log
StandardError=append:/ganuda/logs/sag-unified.log

# Resource limits
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
