# scripts/systemd/sanctuary-state.service
[Unit]
Description=Sanctuary State Daemon - Drift Detection Phase 3A
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=dereadi
WorkingDirectory=/ganuda
ExecStart=/usr/bin/python3 /ganuda/daemons/sanctuary_state.py
Restart=always
RestartSec=60
Environment=PYTHONPATH=/ganuda
EnvironmentFile=/ganuda/config/secrets.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sanctuary-state

[Install]
WantedBy=multi-user.target
