# scripts/systemd/telegram-chief.service
[Unit]
Description=Cherokee AI Telegram Chief Bot
After=network.target llm-gateway.service

[Service]
Type=simple
User=dereadi
Group=dereadi
WorkingDirectory=/ganuda/telegram_bot
Environment=PATH=/home/dereadi/cherokee_venv/bin:/usr/bin:/bin
Environment=PYTHONPATH=/ganuda/lib
EnvironmentFile=/ganuda/config/secrets.env
ExecStart=/home/dereadi/cherokee_venv/bin/python telegram_chief.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-chief

[Install]
WantedBy=multi-user.target
