# scripts/systemd/tribal-vision.service
[Unit]
Description=Tribal Vision - Cherokee AI Camera Intelligence
Documentation=file:///ganuda/services/vision/tribal_vision.py
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=dereadi
Group=dereadi
WorkingDirectory=/ganuda/services/vision
Environment="CHEROKEE_DB_HOST=192.168.132.222"
Environment="CHEROKEE_DB_NAME=zammad_production"
Environment="CHEROKEE_DB_USER=claude"
Environment="CHEROKEE_DB_PASS=jawaseatlasers2"
ExecStart=/home/dereadi/cherokee_venv/bin/python3 /ganuda/services/vision/tribal_vision.py --camera both --interval 60 --model n
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/ganuda/data/vision /tmp
PrivateTmp=false

[Install]
WantedBy=multi-user.target
