# services/consultation-ring.service
[Unit]
Description=Ganuda Consultation Ring — Multi-Model Consultation Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=dereadi
WorkingDirectory=/ganuda
ExecStart=/ganuda/venv/bin/uvicorn services.consultation_ring:app --host 0.0.0.0 --port 9400
Restart=on-failure
RestartSec=5
Environment=PYTHONPATH=/ganuda

[Install]
WantedBy=multi-user.target
