# systemd/executive_jr.service
[Unit]
Description=Cherokee Constitutional AI - Executive Jr (Resource Coordination)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_GROUP
WorkingDirectory=/opt/ganuda_ai
Environment="PATH=/opt/ganuda_ai/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/ganuda_ai/venv/bin/python3 -u /opt/ganuda_ai/daemons/executive_jr_autonomic.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/ganuda_ai/executive_jr.log
StandardError=append:/var/log/ganuda_ai/executive_jr.error.log

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ganuda_ai

[Install]
WantedBy=multi-user.target
