# systemd/memory-jr-autonomic.service
[Unit]
Description=Memory Jr Autonomic Daemon - Cherokee Constitutional AI
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=dereadi
Group=dereadi
WorkingDirectory=/ganuda/daemons
ExecStart=/usr/bin/python3 -u /ganuda/daemons/memory_jr_autonomic.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=memory-jr

# Resource limits (low priority autonomic process)
Nice=10
CPUQuota=10%
MemoryMax=500M

[Install]
WantedBy=multi-user.target
