# config/systemd/secrets-sync.timer
[Unit]
Description=Cherokee Federation — secrets.env push cadence (redfin-only)

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Unit=secrets-sync.service
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
