# config/systemd/secrets-drift.service
[Unit]
Description=Cherokee Federation — detect secrets.env drift across 5 nodes
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=dereadi
Group=dereadi
ExecStart=/ganuda/scripts/secrets-drift-detect.sh
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
