# config/systemd/secrets-drift.timer
[Unit]
Description=Cherokee Federation — secrets.env drift detection cadence (bluefin-only)

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=secrets-drift.service
Persistent=true
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
