# config/systemd/secrets-sync.service
[Unit]
Description=Cherokee Federation — push secrets.env from redfin to other 4 nodes
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dereadi
Group=dereadi
ExecStart=/ganuda/scripts/secrets-sync-push.sh
# Exit code 1 (partial push failure) is reported but not treated as hard fail —
# next timer cycle will retry. Exit code 2 (unreadable source) IS fatal.
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
