# config/systemd/thermal-forget.service
[Unit]
Description=Thermal Forgetting Protocol — DC-9 Memory Pruning
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dereadi
WorkingDirectory=/ganuda
ExecStart=/home/dereadi/cherokee_venv/bin/python3 /ganuda/scripts/thermal_forget.py
Environment=HOME=/home/dereadi
Environment=PYTHONPATH=/ganuda:/ganuda/lib
EnvironmentFile=/ganuda/config/secrets.env
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
